# resources/docker/etc/supervisord.d/build.conf
[program:nginx]
command=su -c "sh /www/build.sh" application
autostart=true
autorestart=false
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
